# .azpipelines/packer/upload/datingapp.service
[Unit]
Description=TFM DatingApp
After=network.target

[Service]
ExecStart=/app/start-app.sh

[Install]
WantedBy=multi-user.target